# servers/mini/systemd/gotify.service
[Unit]
Description=Gotify push notification server (mini)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.config/gotify
ExecStart=%h/.local/bin/gotify
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/gotify %h/.config/gotify

[Install]
WantedBy=default.target
