# servers/mini/systemd/ollama-health.service
[Unit]
Description=Mini LLM stack health check (F21 Phase 3a)
# Don't run until the ollama daemon itself is up.
After=ollama.service
Requires=ollama.service

[Service]
Type=oneshot
# Uses the local fast-path checker (no SSH), not `pl mini llm health`.
# The pl-side command is dev-only and assumes SSH back to mini; running
# it here would self-loop.
ExecStart=%h/.local/bin/ollama-health-check
# Fail loudly if the check script returns non-zero; systemctl --failed
# will pick it up and the journal carries the detail.
StandardOutput=journal
StandardError=journal
# Tag the check's stdout with the unit name so `journalctl -u
# ollama-health.service` shows both the systemd unit events AND the
# check script's OK/FAIL lines in one place. Without this, stdout
# gets journalled under `_COMM=ollama-health-check` and is invisible
# to a unit-scoped query.
SyslogIdentifier=ollama-health.service
