# servers/mini/systemd/quokka-toggle.service
[Unit]
Description=Quokka voice-agent toggle (Gotify watcher)
After=network-online.target gotify.service
Wants=network-online.target

[Service]
Type=simple
# Run via login shell so tmux (and its children like arecord) get the full
# user environment needed for ALSA device access.
ExecStart=/bin/bash -l -c 'exec %h/.local/bin/quokka-toggle'
Restart=on-failure
RestartSec=5

# Light hardening only — voice-agent (launched via tmux) needs full ALSA
# access, real /tmp for scratch WAVs, and read access to ~/. for models.
# The gotify.service-level hardening (PrivateTmp, ProtectHome=read-only)
# breaks arecord/sox when the tmux server inherits the namespace.
NoNewPrivileges=true

[Install]
WantedBy=default.target
