# servers/mayo1/systemd/mayo-fixtures-publish.timer
[Unit]
Description=Nightly Mayo fixture publish (02:00 AEST)

[Timer]
OnCalendar=*-*-* 16:00:00 UTC
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
